# etc/profile.d/history.sh
#!/bin/sh

export HISTFILE=/dev/null
export LESSHISTFILE=/dev/null
